# app-admin/telegraf/files/init/telegraf
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$
CONFIG_FILE=${CONFIG_FILE:="/etc/telegraf/telegraf.conf"}

user=${user:-telegraf}
group=${group:-telegraf}

description="Telegraf Daemon"
command="/usr/bin/telegraf"
command_args=${command_args:="-C $CONFIG_FILE"}
pidfile=${pidfile:-"/var/run/telegraf/telegraf.pid"}
start_stop_daemon_args=" -b -m -u ${user} -g ${group} -- -config ${CONFIG_FILE}"
required_files=$CONFIG_FILE

depend() {
    need net localmount
    use logger
    after bootmisc
}
