# app-admin/influxdb/files/init/influxd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$
CONFIG_FILE=${CONFIG_FILE:="/etc/influxdb/influxdb.conf"}

user=${user:-influxdb}
group=${group:-influxdb}

description="InfluxDB Daemon"
command="/usr/bin/influxd"
command_args=${command_args:="-C $CONFIG_FILE"}
pidfile=${pidfile:-"/var/run/influxdb.pid"}
start_stop_daemon_args=" -b -m -u ${user} -g ${group} -- -config ${CONFIG_FILE}"
required_files=$CONFIG_FILE

depend() {
    need net localmount
    use logger
    after bootmisc
}
